# st/crontab
0 */3 * * * python3.10 /ext/write_only_cern.py
